# crontab
@reboot sleep 30 && /usr/bin/python3 /home/pi/Documents/Testing.py >> /home/pi/>
@reboot sleep 30 && /usr/bin/python3 /home/pi/Documents/CommentAI.py >> /home/p>
